# scripts/proxy-agent.service
[Unit]
Description=Pickle reverse-proxy control agent
Documentation=https://github.com/pnuops/pickle-proxy-agent docs/api/internal.md
After=network-online.target nginx.service
# Wants (not Requires): with Requires, stopping/failing nginx would kill the agent
# and Restart= would not bring it back. The agent degrades gracefully without nginx
# (nginx -t/reload just fail per request), so a soft dependency is correct.
Wants=network-online.target nginx.service

[Service]
Type=simple
# Secrets and tuning come from an env file placed by deploy.sh (mode 0600, root).
# It defines at least PICKLE_PROXY_AGENT_TOKEN; see credentials.md for the source.
EnvironmentFile=/etc/pickle-proxy-agent/agent.env
ExecStart=/usr/local/bin/pickle-proxy-agent
Restart=always
RestartSec=2s
User=root
Group=root

# Hardening. ProtectSystem=full + ReadWritePaths DO work in the unprivileged LXC —
# the sshgw units run exactly that combination in the same environment. An earlier
# 226/NAMESPACE here was misattributed to them; the likely culprits were ProtectHome
# and/or ReadWritePaths entries that did not exist yet. So: no ProtectHome, and the
# certbot paths (absent until certbot is installed) are '-'-prefixed. deploy.sh
# pre-creates every other listed directory.
StateDirectory=pickle-proxy-agent
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/etc/nginx/pickle.d /var/lib/pickle-proxy-agent
ReadWritePaths=-/etc/letsencrypt -/var/www/certbot
# nginx is reloaded via `nginx -s reload` (uid-0 signal to a uid-0 master, no
# capability needed) and certbot/state writes are all root-owned files, so the
# bounding set can be empty.
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
RestrictSUIDSGID=true
# MemoryDenyWriteExecute is deliberately NOT set: its seccomp filter is inherited by
# exec'd children, and certbot (Python/cffi) is not verified to tolerate it.

[Install]
WantedBy=multi-user.target
