# scripts/proxy-agent.service
[Unit]
Description=Pickle reverse-proxy control agent
# Contract: docs/api/internal.md (workspace docs repo)
Documentation=https://github.com/pnuops/pickle-proxy-agent
After=network-online.target nginx.service
# Wants (not Requires): with Requires, stopping/failing nginx would kill the agent
# and Restart= would not bring it back. The agent degrades gracefully without nginx
# (nginx -t/reload just fail per request), so a soft dependency is correct.
Wants=network-online.target nginx.service

[Service]
Type=simple
# Secrets and tuning come from an env file placed by deploy.sh (mode 0600, root).
# It defines at least PICKLE_PROXY_AGENT_TOKEN; see credentials.md for the source.
EnvironmentFile=/etc/pickle-proxy-agent/agent.env
ExecStart=/usr/local/bin/pickle-proxy-agent
Restart=always
RestartSec=2s
User=root
Group=root

# Hardening. ProtectSystem=full + ReadWritePaths DO work in the unprivileged LXC —
# the sshgw units run exactly that combination in the same environment. An earlier
# 226/NAMESPACE here was misattributed to them; the likely culprits were ProtectHome
# and/or ReadWritePaths entries that did not exist yet. So: no ProtectHome, and the
# certbot paths (absent until certbot is installed) are '-'-prefixed. deploy.sh
# pre-creates every other listed directory.
StateDirectory=pickle-proxy-agent
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/etc/nginx/pickle.d /var/lib/pickle-proxy-agent
ReadWritePaths=-/etc/letsencrypt -/var/www/certbot
# The bounding set cannot be empty — both verified live 2026-07-13 by broken applies:
#  - CAP_DAC_OVERRIDE: `nginx -s reload`/-t re-open the error log (www-data:adm 0640
#    on Debian); root without DAC override gets EACCES ("open() error.log failed (13)").
#  - CAP_NET_BIND_SERVICE: `nginx -t` test-binds the listen sockets, and :80/:443 are
#    privileged ports ("bind() to 0.0.0.0:80 failed (13)").
# The uid-0→uid-0 SIGHUP itself needs no capability.
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_NET_BIND_SERVICE
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
RestrictSUIDSGID=true
# MemoryDenyWriteExecute is deliberately NOT set: its seccomp filter is inherited by
# exec'd children, and certbot (Python/cffi) is not verified to tolerate it.

[Install]
WantedBy=multi-user.target
